# rag-qdrant-local/installer/systemd/reineke-rag-backup.service
[Unit]
Description=Reineke-RAG nightly backup
After=reineke-rag.service
Requires=reineke-rag.service

[Service]
Type=oneshot
WorkingDirectory=/opt/reineke-rag
ExecStart=/opt/reineke-rag/scripts/backup.sh
StandardOutput=journal
StandardError=journal
